# etc/service/phpfpm/run
#!/usr/bin/env bash

echo 'Starting phpfpm'

php-fpm${PHP_VERSION} -F -O
