# etc/service/phpfpm/run
#!/usr/bin/env bash

php-fpm${PHP_VERSION} -F -O