# etc/service/websocket/run
#!/usr/bin/env bash

php ${WORK_DIR}/artisan websockets:serve