# etc/service/websocket/run
#!/usr/bin/env bash

echo 'Starting Websockets'

php ${WORK_DIR}/artisan websockets:serve
