# etc/service/horizon/run
#!/usr/bin/env bash

echo 'Starting Horizon'

php ${WORK_DIR}/artisan horizon
