# etc/service/nginx/run
#!/usr/bin/env bash

echo 'Starting Nginx'

nginx -g 'daemon off;'
